# src/backend/etc/supervisord.ini
[supervisord]
logfile = /app/supervisord.log
nodaemon = true

[program:aiohttp]
numprocs = 4
numprocs_start = 1
process_name = aiohttp_%(process_num)s
redirect_stderr = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

command=aisauna --config /app/etc/development.yml --port=808%(process_num)s
              --socket_path=/var/lib/aiohttp/%(program_name)s_%(process_num)s.sock

autostart=true
autorestart=true
